# project/config/production/crontab/mvc_frame
* * * * * root /bin/bash /var/www/mvc_frame/project/tool/production/check_update.sh
